# scripts/com.alphehelix.xbot.synthesis.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  每天三次（台灣時間 8:00 / 15:00 / 20:30）對累積的所有作者貼文做跨作者觀點彙整，
  並更新網站、自動 push、寄信。排程時間依 Mac 系統本地時區，請確認系統時區為 Asia/Taipei。
  安裝方式：
    1. 把 __PYTHON__ 換成 venv 內 python 絕對路徑、__PROJECT_DIR__ 換成專案絕對路徑。
    2. cp scripts/com.alphehelix.xbot.synthesis.plist ~/Library/LaunchAgents/
    3. launchctl load ~/Library/LaunchAgents/com.alphehelix.xbot.synthesis.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alphehelix.xbot.synthesis</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>-m</string>
        <string>src.main</string>
        <string>synthesis</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key><integer>8</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>15</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>20</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>__PROJECT_DIR__/xbot.log</string>
    <key>StandardErrorPath</key>
    <string>__PROJECT_DIR__/xbot.log</string>
</dict>
</plist>
